# dockerbuild/supervisord.conf
[supervisord]
nodaemon=false

[program:mongod]
command=mongod
