# dockerbuild/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod
